# templates/plone_upstart.conf.erb
# Plone - Web-Content Management System
#
# Based on Python and ZOPE

description "start plone instance <%= @name %>"
author "Fabio Rauber based on previous work by Josh Sehn and Christoph Glaubitz"
version "0.4"

console none
respawn

limit nofile 65535 65535

start on (local-filesystems and net-device-up and runlevel [2345])
stop on runlevel [!2345]

post-stop exec sleep 5

exec <%= @install_dir %>/<%= @name %>/bin/instance console
